# spore/systemd/myconex.service
[Unit]
Description=MYCONEX Discord AI Mesh Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=techno-shaman
Group=techno-shaman
WorkingDirectory=/home/techno-shaman/myconex
Environment=PYTHONPATH=/home/techno-shaman/myconex
EnvironmentFile=/home/techno-shaman/myconex/.env
ExecStart=/usr/bin/python3 /home/techno-shaman/myconex/main.py --mode full
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30

# Logging → journalctl -u myconex -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=myconex

[Install]
WantedBy=multi-user.target